# arm64/config/includes.chroot_after_packages/etc/systemd/system/pai-setup.service
[Unit]
Description=PAI First-Boot Model Setup
After=ollama.service network-online.target
Wants=ollama.service network-online.target
ConditionPathExists=!/tmp/.pai-setup-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pai-setup
ExecStartPost=/usr/bin/touch /tmp/.pai-setup-done
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
